# noip2.conf
description "only execute noip2 when ethernet port 0 is workable"

start on filesystem and net-device-up IFACE=eth0
stop on runlevel [016]

exec /usr/local/bin/noip2
